# supervisord.conf
[program:worker]
command = django-admin.py celery worker --concurrency=2 --beat --events
directory = /home/dotcloud/current/