# supervisord.conf
[program:worker]
command = django-admin.py celery worker --concurrency=2 --beat --events -n worker-main
directory = /home/dotcloud/current/
numprocs = 1